# deploy/sish-retrain.service
[Unit]
Description=SISH weekly full retrain and recalibration
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=buoycast
WorkingDirectory=/opt/sish
ExecStart=/bin/bash -c 'set -e; V=/opt/sish/venv/bin/python; $V fetch.py; $V fetch_weather.py; $V train_q.py --refit-full; $V backtest.py; $V corr.py; $V publish.py'
TimeoutStartSec=7200
Nice=15
