# deploy/sish-retrain.timer
[Unit]
Description=run SISH retrain weekly

[Timer]
OnCalendar=Sun *-*-* 09:30:00 UTC
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
